# deploy/systemd/pantry-manager.service
[Unit]
Description=Pantry Manager FastAPI
After=network.target

[Service]
User=pantry
Group=pantry
WorkingDirectory=/opt/pantry-manager
Environment=APP_VERSION=0.1.0
Environment=DEPLOY_ENV=production
Environment=GIT_SHA=unknown
ExecStart=/opt/pantry-manager/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
